# latex-compiler.service
[Unit]
Description=LaTeX Compiler Service
After=network.target
Wants=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/latex-compiler
Environment=API_KEY=djsakjc213hjbkk3h123jkb123kbj
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/latex-compiler
ExecStart=/usr/local/bin/gunicorn --config /opt/latex-compiler/gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=latex-compiler

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/latex-compiler/logs /var/log/gunicorn
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
